# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=-/sbin/agetty --autologin imphnen --noclear %I 38400 $TERM
#ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root - $TERM
